# tilemill/socat.conf
[program:tilemill]
user=root
command=socat TCP-LISTEN:5432,reuseaddr,fork TCP:$PG_PORT_5432_TCP_ADDR:$PG_PORT_5432_TCP_PORT
autorestart=true
stopsignal=INT

